# environments/mezzapp/mezzapp.conf
{% set env_base = pillar.get('environment', {}) %}
{% set env = env_base.get('mezzapp', {}) %}

{% set cmd_args = [
    '--http-socket :3000',
    '--pidfile=wsgi.pid',
    '--workers {}'.format(grains.get('mezzapp_web_workers', 7)),
    '--master',
    '--stats 0.0.0.0:23000',
    '--wsgi mezzapp_wsgi',
    '--max-requests 10',
    '--harakiri 300',
    '--static-map /static=/home/mezzapp/static_root/static',
] %}

{% if env.get('media_files', {}).get('filesystem_root', False) %}
# Static map media files
{% do cmd_args.append("--static-map /media={}".format(env['media_files']['filesystem_root'])) %}
{% else %}
# No media files static map
{% endif %}


[program:mezzapp]

command=/home/mezzapp/vmezzapp/bin/uwsgi {{' '.join(cmd_args)}}

stdout_logfile=logs/wsgi.log
redirect_stderr=true
stopsignal=INT
stopasgroup=true
killasgroup=true
environment=MEZZAPP_LOG_LOCATION=/home/mezzapp/logs,DJANGO_SETTINGS_MODULE=local_settings,PYTHONPATH=/home/mezzapp/config

