# environments/mezzapp/mezzapp_sup.conf

[program:mezzapp_sup]
command=/usr/bin/supervisord
user=mezzapp
directory=/home/mezzapp
stdout_logfile=/var/log/supervisor/mezzapp_sup.log
redirect_stderr=true

